# restic-backup.timer
[Unit]
Description=Run Restic Backup 5 minutes after boot

[Timer]
OnBootSec=5min
Unit=restic-backup.service

[Install]
WantedBy=timers.target
